# files/docker/docker-prune.timer

[Unit]
Description=Daily docker system prune

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
